# backend/cron/crontab
00 4 * * * ~/wordteacher/wordteacher_renew_cert.sh
00 4 * * * ~/wordteacher/mongo_backup.sh
