# backend/cron/crontab
00 4 * * * ~/wordteacher/wordteacher_renew_cert.sh
